# roles/rabbitmq/templates/rabbitmq-server.service.j2
[Service]
ExecStop=
ExecStop=/usr/lib/rabbitmq/bin/rabbitmqctl shutdown
